# meshtastic-dashboard.service
[Unit]
Description=Meshtastic Network Dashboard
Documentation=https://github.com/apocas/meshtastic-dashboard
After=network.target
Wants=network.target

[Service]
Type=notify
User=meshtastic
Group=meshtastic
WorkingDirectory=/home/meshtastic/meshtastic-dashboard
Environment=PATH=/home/meshtastic/meshtastic-dashboard.venv/bin
Environment=PYTHONPATH=/home/meshtastic/meshtastic-dashboard
Environment=FLASK_ENV=production
Environment=FLASK_SECRET_KEY=your-secret-key-change-this-in-production
ExecStart=/home/meshtastic/meshtastic-dashboard.venv/bin/gunicorn --config gunicorn.conf.py wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=meshtastic-dashboard

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/meshtastic/meshtastic-dashboard
ReadWritePaths=/var/log/meshtastic-dashboard
ReadWritePaths=/var/run/meshtastic-dashboard

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Graceful shutdown
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
